# config/LaunchAgents/com.zen-mcp-server.logrotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zen-mcp-server.logrotate</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/zen-logrotate</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>__HOME__/.local/log/zen-logrotate.log</string>
    
    <key>StandardErrorPath</key>
    <string>__HOME__/.local/log/zen-logrotate.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:__HOME__/.local/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>USER</key>
        <string>__USER__</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist> 